# pi_controller_boy/pishutdown.service
[Service]
ExecStart=/usr/bin/python /home/pi/rpi_code_jessie/pi_controller_boy/pishutdown.py
WorkingDirectory=/home/pi/rpi_code_jessie/pi_controller_boy/
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=pishutdown
User=root
Group=root

[Install]
WantedBy=multi-user.target
