# jessie_temp_logger/templogger.service
[Unit]
Description=TempLogger

[Service]
Type=simple
ExecStart=/home/pi/rpi_code_jessie/jessie_temp_logger/templogger.py
